# ops/systemd/sakiot-web.service
[Unit]
Description=Sakiot web server
After=network-online.target postgresql.service
Wants=network-online.target

[Service]
Type=simple
User=sakiot
Group=sakiot
WorkingDirectory=/var/lib/sakiot/data
EnvironmentFile=/etc/sakiot/production.env
EnvironmentFile=/srv/sakiot/current/web/service.env
ExecStart=/srv/sakiot/current/web/web_server
Restart=on-failure
RestartSec=2
TimeoutStopSec=30
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true

[Install]
WantedBy=multi-user.target
