# ops/systemd/sakiot-db-backup-nightly.timer
[Unit]
Description=Run the nightly Sakiot database backup

[Timer]
OnCalendar=*-*-* 03:17:00
Persistent=true
AccuracySec=1m
Unit=sakiot-db-backup@nightly.service

[Install]
WantedBy=timers.target
